# systemd/snapper-timeline.service.d/10-snappersend.conf
# Drop-in on snapper-timeline.service — the snappersend trigger.
#
# Installed to /etc/systemd/system/snapper-timeline.service.d/10-snappersend.conf
# (NOT by editing the packaged /usr/lib/systemd/system unit, which is clobbered
# on upgrade). When the timeline timer pulls in the snapshot service each hour,
# this also pulls in snappersend, ordered strictly AFTER the snapshot completes
# so the new snapshot is already on disk when snappersend picks the newest.
#
# Wants= (not Requires=): snappersend still runs on the hours snapper's
# --timeline no-ops (it just decides there's nothing new to send and exits 0),
# and a snappersend failure — e.g. a transient network blip to the destination,
# which legitimately exits non-zero — never propagates back to mark
# snapper-timeline.service itself as failed.
[Unit]
Wants=snappersend.service
Before=snappersend.service
