# systemd/snapper-timeline.service.d/20-bootmirror.conf
# Drop-in on snapper-timeline.service — the bootmirror-sync trigger.
#
# Installed to /etc/systemd/system/snapper-timeline.service.d/20-bootmirror.conf,
# alongside 10-snappersend.conf and following the same pattern: each timeline
# tick also pulls in bootmirror-sync.service, which orders itself Before= the
# timeline (its own unit file), so snapper snapshots a freshly-synced
# @bootmirror in the same tick as root/home.
#
# Wants= (not Requires=): a soft dependency only. snapper does not wait on the
# sync's exit status and a sync failure can never fail or block the timeline —
# the mirror is then simply one tick stale.
[Unit]
Wants=bootmirror-sync.service
