# systemd/snapper-timeline.service.d/30-manifest.conf
# Drop-in on snapper-timeline.service — the restore-manifest trigger.
#
# Installed to /etc/systemd/system/snapper-timeline.service.d/30-manifest.conf,
# alongside 10-snappersend.conf and 20-bootmirror.conf and following the same
# pattern: each timeline tick also pulls in snappersend-manifest.service, which
# orders itself Before= the timeline (its own unit file), so snapper snapshots
# a fresh restore manifest in the same tick — the manifest the destination
# receives describes the very snapshot that carries it.
#
# Wants= (not Requires=): a soft dependency only. snapper does not wait on the
# write's exit status and a manifest failure can never fail or block the
# timeline — a restore then falls back to the last successfully written one.
[Unit]
Wants=snappersend-manifest.service
