# systemd/snapsend-watchdog.timer
[Unit]
Description=Periodic staleness check for di-snapsend replication
Documentation=file:/usr/local/share/doc/di-snapsend/README.md

[Timer]
# A few times a day is plenty to catch a stalled replication chain.
OnCalendar=*-*-* 0/6:17:00
RandomizedDelaySec=600
Persistent=true

[Install]
WantedBy=timers.target
